# rootfs/etc/services.d/mosquitto/run
#!/usr/bin/with-contenv sh

/usr/sbin/mosquitto -c $MOSQUITTO_CONF_DIR/solar.conf
