# rootfs/etc/services.d/mosquitto/run
#!/usr/bin/with-contenv sh

/usr/sbin/mosquitto -c /etc/mosquitto/mosquitto.conf
